# ops/chimera-api.service
[Unit]
Description=Chimera LogMind API (UDS)
After=network.target

[Service]
Type=simple
User=chimera
Group=chimera
SupplementaryGroups=systemd-journal
Environment=CHIMERA_API_SOCKET=/run/chimera/api.sock
Environment=CHIMERA_DB_PATH=/var/lib/chimera/chimera.duckdb
WorkingDirectory=/opt/chimera
ExecStart=/opt/chimera/venv/bin/python3 -m api.server
RuntimeDirectory=chimera
RuntimeDirectoryMode=0750
StateDirectory=chimera
Restart=on-failure

[Install]
WantedBy=multi-user.target
